# dashboardinstall/systemd/fastsync-localhost-check.service
[Unit]
Description=Check the status of fastsync localhost

[Service]
ExecStart=/etc/monitor-scripts/fastsync-localhost.sh
